# scripts/com.researchlab.brief-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Actionable Market Brief — timer launchd job (this-MacBook host; template).
  Runs the timer wrapper scripts/brief-refresh-and-push.sh at the four ET windows:
  it refreshes the deterministic data (brief-refresh.mjs), REGENERATES the narrative
  (market-brief.payload.json) with the GitHub Copilot CLI (Opus 4.8, file-edits only),
  COMMITS the changed brief files, and git-pushes so GitHub Pages redeploys. Push auth =
  the repo's HTTPS remote + macOS osxkeychain helper; the Copilot CLI reuses its own login
  (both headless while you are logged in). evo-x2 + a knb systemd timer is an equivalent
  always-on alternative host.

  INSTALL (edit the ABSOLUTE path to the wrapper below to match your checkout, then):
    cp scripts/com.researchlab.brief-refresh.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.researchlab.brief-refresh.plist

  Runs at the four ET windows. launchd fires in LOCAL time, so the hours below assume
  the machine is on America/New_York; adjust if your box is on another zone. The wrapper
  auto-detects the window from the ET clock and no-ops on weekends (a weekend fire simply
  commits nothing), so no per-run args are needed.
-->
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.researchlab.brief-refresh</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>/ABSOLUTE/PATH/TO/research-lab/scripts/brief-refresh-and-push.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <array>
      <dict><key>Hour</key><integer>7</integer><key>Minute</key><integer>30</integer></dict>
      <dict><key>Hour</key><integer>11</integer><key>Minute</key><integer>0</integer></dict>
      <dict><key>Hour</key><integer>15</integer><key>Minute</key><integer>0</integer></dict>
      <dict><key>Hour</key><integer>17</integer><key>Minute</key><integer>0</integer></dict>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/opt/homebrew/bin:/opt/local/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/tmp/brief-refresh.out.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/brief-refresh.err.log</string>
  </dict>
</plist>
