# scripts/com.researchlab.brief-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Actionable Market Brief — isolated launchd job (this-MacBook host; template).
  Runs scripts/brief-refresh-scheduled.sh at the four ET windows. The launcher clones
  origin/main into a disposable checkout, then invokes brief-refresh-and-push.sh to
  refresh the deterministic data (brief-refresh.mjs), REGENERATE the narrative
  (market-brief.payload.json) with the GitHub Copilot CLI (Opus 4.8, file-edits only),
  COMMITS the changed brief files, and git-pushes so GitHub Pages redeploys. Push auth =
  the repo's HTTPS remote + macOS osxkeychain helper; the Copilot CLI reuses its own login
  (both headless while you are logged in). evo-x2 + a knb systemd timer is an equivalent
  always-on alternative host.

  INSTALL (edit the ABSOLUTE path to the wrapper below to match your checkout, then):
    cp scripts/com.researchlab.brief-refresh.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.researchlab.brief-refresh.plist

  Runs at the four ET windows. launchd fires in LOCAL time and cannot be given a timezone,
  so the hours below are the LOCAL (US Pacific, ET-3) equivalents of the ET windows — NOT
  the raw ET hours (04:30/08:00/12:00/14:00 PT = 07:30/11:00/15:00/17:00 ET). If this box
  moves off US Pacific, re-derive them (local = ET - 3). The wrapper auto-detects the window
  from the ET clock and no-ops on weekends (a weekend fire simply commits nothing), so no
  per-run args are needed.
-->
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.researchlab.brief-refresh</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>/ABSOLUTE/PATH/TO/research-lab/scripts/brief-refresh-scheduled.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <array>
      <dict><key>Hour</key><integer>4</integer><key>Minute</key><integer>30</integer></dict> <!-- 07:30 ET pre-market -->
      <dict><key>Hour</key><integer>8</integer><key>Minute</key><integer>0</integer></dict>  <!-- 11:00 ET morning -->
      <dict><key>Hour</key><integer>12</integer><key>Minute</key><integer>0</integer></dict> <!-- 15:00 ET pre-close -->
      <dict><key>Hour</key><integer>14</integer><key>Minute</key><integer>0</integer></dict> <!-- 17:00 ET after-hours -->
    </array>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/opt/homebrew/bin:/opt/local/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/tmp/brief-refresh.out.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/brief-refresh.err.log</string>
  </dict>
</plist>
